# utils/systemd.service/nichi.service
[Unit]
Description=sigcoww pos

[Service]
ExecStart=/usr/local/bin/node /home/pi/nichi/index.js
Restart=always

[Install]
WantedBy=multi-user.target
